# scripts/deployment/linux/pocketportal.service
[Unit]
Description=PocketPortal AI Agent
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/pocketportal
Environment="PATH=/home/%i/pocketportal/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/%i/pocketportal/venv/bin/pocketportal start --interface telegram
Restart=on-failure
RestartSec=30s
StandardOutput=append:/home/%i/pocketportal/logs/systemd.log
StandardError=append:/home/%i/pocketportal/logs/systemd_error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/pocketportal/logs /home/%i/pocketportal/data

[Install]
WantedBy=multi-user.target
